# deploy/openclaw.service
[Unit]
Description=OpenClaw Deal Scout Gateway (MCP orchestrator + HTTP server)
After=network.target

[Service]
Type=simple
User=asmaiqbal01
WorkingDirectory=/home/asmaiqbal01/openclaw-deal-scout
ExecStart=/usr/bin/python3.12 -m openclaw_gateway
EnvironmentFile=/home/asmaiqbal01/openclaw-deal-scout/.env
Environment=SCHEDULER_MODE=gateway

[Install]
WantedBy=multi-user.target
