# deploy/openclaw@.service
[Unit]
Description=OpenClaw Deal Scout Pipeline (one-shot cycle)
After=network.target

[Service]
Type=oneshot
User=%i
WorkingDirectory=/home/%i/openclaw-deal-scout
ExecStart=/usr/bin/python3.12 -m pipeline_orchestrator
EnvironmentFile=/home/%i/openclaw-deal-scout/.env
Environment=SCHEDULER_MODE=systemd

[Install]
WantedBy=multi-user.target
